# root/etc/services.d/xeoma/run
#!/usr/bin/with-contenv bash

exec s6-setuidgid dockeruser $INSTALL_LOCATION/xeoma -core -log -startdelay 5 
